# overlay/system/etc/init/halium-debug.conf
# Halium emergency telnet - starts BEFORE android mounts
# Allows debug access even if Android container never starts
# Connect: telnet 10.15.19.82 23 (USB tethering must be active)

description "Halium emergency telnet"

# Start as early as possible - only needs network
start on (net-device-up IFACE=rndis0
    or net-device-up IFACE=usb0)

stop on runlevel [016]

respawn
respawn limit unlimited 0

exec /usr/sbin/telnetd -l /bin/bash -F -p 23
