# overlay/system/etc/init/mount-android.conf
# Mount device partitions for Halium/Ubuntu Touch
# Samsung Galaxy Tab A7 Lite (gta7litewifi) - MT8768 - Halium 12

description "Mount device partitions"

start on (startup
    and started udev)

emits mounted-android

task

script
    # --- cpuset cgroup ---
    mkdir -p /dev/cpuset
    if ! mountpoint -q /dev/cpuset; then
        mount none /dev/cpuset -t cpuset -o nodev,noexec,nosuid || true
    fi
    # Create expected cpuset groups
    mkdir -p /dev/cpuset/foreground
    mkdir -p /dev/cpuset/background
    mkdir -p /dev/cpuset/system-background
    mkdir -p /dev/cpuset/top-app
    mkdir -p /dev/cpuset/restricted

    # --- binderfs for Android 12 ---
    if ! mountpoint -q /dev/binderfs; then
        mkdir -p /dev/binderfs
        mount -t binder binder /dev/binderfs \
            -o stats=global,lazy 2>/dev/null || true
    fi
    # Symlink legacy binder nodes for libhybris compatibility
    for node in binder hwbinder vndbinder; do
        if [ ! -e /dev/$node ] && [ -e /dev/binderfs/$node ]; then
            ln -sf /dev/binderfs/$node /dev/$node 2>/dev/null || true
        fi
    done

    # --- Run main partition mount script ---
    /usr/sbin/mount-android.sh

    # --- Prevent gralloc.default.so fallback (kills PowerVR GE8320) ---
    if [ -e /android/vendor/lib64/hw/gralloc.default.so ]; then
        mount -o bind /dev/null /android/vendor/lib64/hw/gralloc.default.so || true
    fi
    if [ -e /android/vendor/lib/hw/gralloc.default.so ]; then
        mount -o bind /dev/null /android/vendor/lib/hw/gralloc.default.so || true
    fi

    # --- Overlay vndservicemanager binary ---
    if [ -f /usr/share/halium-overlay/vendor/bin/vndservicemanager ]; then
        mount -o bind /usr/share/halium-overlay/vendor/bin/vndservicemanager \
            /android/vendor/bin/vndservicemanager || true
    fi

    # --- Overlay vndservicemanager RC (disables time_daemon conflicts) ---
    if [ -f /usr/share/halium-overlay/vendor/etc/init/vndservicemanager.rc ]; then
        mkdir -p /android/vendor/etc/init
        mount -o bind /usr/share/halium-overlay/vendor/etc/init/vndservicemanager.rc \
            /android/vendor/etc/init/vndservicemanager.rc || true
    fi

    # --- Overlay disable-bt RC ---
    if [ -f /usr/share/halium-overlay/vendor/etc/init/disable-bt.rc ]; then
        mkdir -p /android/vendor/etc/init
        mount -o bind /usr/share/halium-overlay/vendor/etc/init/disable-bt.rc \
            /android/vendor/etc/init/disable-bt.rc || true
    fi

    # --- Overlay USB RC (MT8768 specific) ---
    if [ -f /usr/share/halium-overlay/vendor/etc/init/hw/init.mt6765.usb.rc ]; then
        mkdir -p /android/vendor/etc/init/hw
        mount -o bind /usr/share/halium-overlay/vendor/etc/init/hw/init.mt6765.usb.rc \
            /android/vendor/etc/init/hw/init.mt6765.usb.rc || true
    fi

    # --- Overlay audio_policy_configuration.xml ---
    if [ -f /usr/share/halium-overlay/vendor/etc/init/audio_policy_configuration.xml ]; then
        # Try common vendor audio config paths on MTK
        for apc in \
            /android/vendor/etc/audio_policy_configuration.xml \
            /android/vendor/etc/init/audio_policy_configuration.xml; do
            if [ -f "$apc" ]; then
                mount -o bind \
                    /usr/share/halium-overlay/vendor/etc/init/audio_policy_configuration.xml \
                    "$apc" || true
            fi
        done
    fi

    # --- Overlay null apexd (prevent APEX mount loops on Halium) ---
    if [ -f /usr/share/halium-overlay/system/bin/apexd ]; then
        if [ -f /android/system/bin/apexd ]; then
            mount -o bind /usr/share/halium-overlay/system/bin/apexd \
                /android/system/bin/apexd || true
        fi
    fi

    # --- Writable dirs for usb_moded (rootfs is read-only squashfs) ---
    # usb_moded tries to write /etc/usb-moded/usb-moded.ini and /etc/modprobe.d/g_ether.conf
    mkdir -p /userdata/usb-moded /userdata/modprobe.d
    if [ ! -f /userdata/modprobe.d/g_ether.conf ]; then
        touch /userdata/modprobe.d/g_ether.conf
    fi
    mount -o bind /userdata/usb-moded /etc/usb-moded || true
    mount -o bind /userdata/modprobe.d /etc/modprobe.d || true

    # --- Cache partition redirect to userdata ---
    umount /android/cache 2>/dev/null || true
    mkdir -p /userdata/cache
    mount -o bind /userdata/cache /android/cache || true
    mkdir -p /android/cache/recovery

    # --- Metadata bind into LXC rootfs ---
    if [ -d /android/metadata ]; then
        mkdir -p /var/lib/lxc/android/rootfs/metadata
        mount -o bind /android/metadata /var/lib/lxc/android/rootfs/metadata || true
    fi

    initctl emit -n mounted-android
end script
