# overlay/system/etc/init/usb-early-debug.conf
# Early USB RNDIS setup for Halium debug
# Must run before mount-android so telnet is reachable during mount failures
# Samsung SM-T225: VID=04E8, RNDIS PID=6863

description "Early USB RNDIS for Halium debug"

start on startup
stop on runlevel [016]

task

script
    # Wait for configfs
    sleep 2
    
    GADGET=/sys/kernel/config/usb_gadget
    
    # Mount configfs if not mounted
    mountpoint -q /sys/kernel/config || mount -t configfs none /sys/kernel/config
    
    [ -d "$GADGET" ] || exit 0  # No configfs gadget support → skip
    
    # Create gadget if not already done by Android
    if [ ! -d "$GADGET/g1" ]; then
        mkdir -p "$GADGET/g1"
        echo "0x04E8" > "$GADGET/g1/idVendor"
        echo "0x6863" > "$GADGET/g1/idProduct"
        mkdir -p "$GADGET/g1/strings/0x409"
        echo "Ubuntu Touch Debug" > "$GADGET/g1/strings/0x409/product"
        echo "Samsung" > "$GADGET/g1/strings/0x409/manufacturer"
        
        mkdir -p "$GADGET/g1/functions/rndis.usb0"
        mkdir -p "$GADGET/g1/configs/c.1"
        mkdir -p "$GADGET/g1/configs/c.1/strings/0x409"
        echo "rndis" > "$GADGET/g1/configs/c.1/strings/0x409/configuration"
        ln -sf "$GADGET/g1/functions/rndis.usb0" "$GADGET/g1/configs/c.1/"
        
        UDC=$(ls /sys/class/udc | head -1)
        [ -n "$UDC" ] && echo "$UDC" > "$GADGET/g1/UDC"
    fi
    
    # Bring up RNDIS interface
    sleep 1
    for iface in rndis0 usb0; do
        if ip link show $iface > /dev/null 2>&1; then
            ip link set $iface up
            ip addr add 10.15.19.82/24 dev $iface 2>/dev/null || true
            break
        fi
    done
end script
