# overlay/system/usr/share/upstart/sessions/mtp-server.conf
description "MTP protocol server"

start on :sys:android-usb-connected or android-usb-connected or :sys:android-mtp-on or android-mtp-on or started unity8
stop on :sys:android-usb-disconnected or android-usb-disconnected or :sys:android-mtp-off or android-mtp-off or desktop-end

pre-start script
    [ "$USER" != "lightdm" ] || { stop; exit 0; }

    # Support both legacy android_usb and configfs gadget
    if [ -e /sys/class/android_usb/android0/state ]; then
        sleep 0.5
        grep -q DISCONNECTED /sys/class/android_usb/android0/state && { stop; exit 0; }
    elif [ -e /sys/kernel/config/usb_gadget/g1/UDC ]; then
        # configfs mode - check if UDC is bound
        UDC_VAL=$(cat /sys/kernel/config/usb_gadget/g1/UDC 2>/dev/null)
        [ -z "$UDC_VAL" ] && { stop; exit 0; }
    else
        stop; exit 0
    fi

    # Check MTP function is active
    if [ ! -e /sys/kernel/config/usb_gadget/g1/configs/c.1/mtp.mtp ] && \
       [ ! -e /sys/class/android_usb/android0/functions ]; then
        stop; exit 0
    fi

    # Fix for LP: #1389223 - UnityGreeter availability race
    if ! dbus-send --session --print-reply \
            --dest=com.canonical.UnityGreeter \
            / org.freedesktop.DBus.Introspectable.Introspect 2>&1 | \
                grep -qc IsActive; then
        sleep 10
    fi
end script

env GOOGLE_LOGTOSTDERR=1

respawn
respawn limit unlimited

exec /usr/bin/mtp-server
