# overlay/system/usr/lib/systemd/system/sensorfwd.service.d/wait-android.conf
[Unit]
# Sensor HAL lives in Android container — must wait for LXC to be ready
# Fixes: HYBRIS CTL setDelay -> -1 (HAL not ready when sensorfwd first probes)
After=lxc@android.service android-mount.service

[Service]
# Retry up to 5 times — HAL may not be ready immediately after container start
Restart=on-failure
RestartSec=5s
StartLimitBurst=5
StartLimitIntervalSec=60s
