# overlay/system/usr/lib/systemd/system/repowerd.service.d/wait-logind.conf
[Unit]
# Wait until logind has processed all sessions before repowerd starts
# Fixes: NoSessionForPID: PID does not belong to any known session
After=systemd-logind.service dbus.service
Requires=systemd-logind.service
